# stage3/00-install-packages/files/strut.service
[Unit]
Description=STRUT

[Service]
ExecStart=/usr/bin/strut
Restart=always

[Install]
WantedBy=multi-user.target
